# services/obsidian-runner/obsidian-runner.service
[Unit]
Description=Obsidian Scheduled Prompt Runner
After=network.target obsidian-sync.service
Wants=obsidian-sync.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/obsidian-runner/obsidian_runner.py
Environment=VAULT_BASE=__VAULT_PATH__
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
